# root/etc/periodic/daily/letsencrypt
#!/bin/sh

/getssl.sh -q -a >> /var/log/cron.log 2>&1
